# WITH_WORKER/worker.service
[Unit]
Description=Search Worker Service
After=network.target

[Service]
Type=simple
User=your_user
Environment=DATABASE_URL=postgresql://user:pass@host:port/db
WorkingDirectory=/path/to/your/app
ExecStart=/usr/bin/python3 /path/to/your/app/WITH_WORKER/worker.py
Restart=always

[Install]
WantedBy=multi-user.target 